# files/agent_groups/linux-auditd-rules.conf
# Linux Auditd Rules for Wazuh Agent Groups
# Deploy to Linux agents as /etc/audit/rules.d/wazuh.rules
#
# These rules generate events that match the auditd decoders and
# the 200110-auditd.xml rule file for MITRE ATT&CK detection.
#
# After deploying, reload: sudo augenrules --load

# Remove existing rules
-D

# Set buffer size (increase for busy systems)
-b 8192

# Failure mode: 1 = printk, 2 = panic
-f 1

# ═══════════════════════════════════════════════════
# Sensitive File Access (T1003, T1552)
# ═══════════════════════════════════════════════════
-w /etc/shadow -p wa -k access
-w /etc/gshadow -p wa -k access
-w /etc/passwd -p wa -k identity
-w /etc/group -p wa -k identity

# ═══════════════════════════════════════════════════
# User and Group Changes (T1136, T1098)
# ═══════════════════════════════════════════════════
-w /usr/sbin/useradd -p x -k identity
-w /usr/sbin/userdel -p x -k identity
-w /usr/sbin/usermod -p x -k identity
-w /usr/sbin/groupadd -p x -k identity
-w /usr/sbin/groupdel -p x -k identity
-w /usr/sbin/groupmod -p x -k identity
-w /usr/bin/passwd -p x -k identity

# ═══════════════════════════════════════════════════
# Privilege Escalation (T1548)
# ═══════════════════════════════════════════════════
-w /etc/sudoers -p wa -k priv_esc
-w /etc/sudoers.d/ -p wa -k priv_esc
-a always,exit -F arch=b64 -S setuid -S setgid -S setreuid -S setregid -k priv_esc

# ═══════════════════════════════════════════════════
# Kernel Module Loading (T1547.006)
# ═══════════════════════════════════════════════════
-w /sbin/insmod -p x -k kernel_module
-w /sbin/modprobe -p x -k kernel_module
-w /sbin/rmmod -p x -k kernel_module
-a always,exit -F arch=b64 -S init_module -S finit_module -S delete_module -k kernel_module

# ═══════════════════════════════════════════════════
# Persistence (T1053, T1543)
# ═══════════════════════════════════════════════════
-w /etc/cron.d/ -p wa -k persistence
-w /etc/crontab -p wa -k persistence
-w /var/spool/cron/ -p wa -k persistence
-w /etc/systemd/system/ -p wa -k persistence
-w /usr/lib/systemd/system/ -p wa -k persistence
-w /etc/init.d/ -p wa -k persistence

# ═══════════════════════════════════════════════════
# SSH Configuration (T1098.004)
# ═══════════════════════════════════════════════════
-w /etc/ssh/sshd_config -p wa -k access
-w /root/.ssh/ -p wa -k access

# ═══════════════════════════════════════════════════
# Shell Configuration (T1546.004)
# ═══════════════════════════════════════════════════
-w /etc/profile -p wa -k persistence
-w /etc/profile.d/ -p wa -k persistence
-w /etc/bash.bashrc -p wa -k persistence
-w /etc/environment -p wa -k persistence
-w /etc/ld.so.preload -p wa -k persistence

# ═══════════════════════════════════════════════════
# Suspicious Process Execution (T1059)
# ═══════════════════════════════════════════════════
# Execution from /tmp, /var/tmp, /dev/shm
-a always,exit -F arch=b64 -S execve -F dir=/tmp -k susp_activity
-a always,exit -F arch=b64 -S execve -F dir=/var/tmp -k susp_activity
-a always,exit -F arch=b64 -S execve -F dir=/dev/shm -k susp_activity

# ═══════════════════════════════════════════════════
# Network Connections (T1071)
# ═══════════════════════════════════════════════════
-a always,exit -F arch=b64 -S connect -F a2=16 -k network_connect

# ═══════════════════════════════════════════════════
# Time Changes (T1070.006)
# ═══════════════════════════════════════════════════
-a always,exit -F arch=b64 -S adjtimex -S settimeofday -S clock_settime -k time_change

# ═══════════════════════════════════════════════════
# File Permission Changes (T1222)
# ═══════════════════════════════════════════════════
-a always,exit -F arch=b64 -S chmod -S fchmod -S fchmodat -k perm_mod
-a always,exit -F arch=b64 -S chown -S fchown -S lchown -S fchownat -k perm_mod

# ═══════════════════════════════════════════════════
# File Deletion (T1070.004)
# ═══════════════════════════════════════════════════
-a always,exit -F arch=b64 -S unlink -S unlinkat -S rename -S renameat -k delete

# ═══════════════════════════════════════════════════
# Log Tampering (T1070.002)
# ═══════════════════════════════════════════════════
-w /var/log/ -p wa -k access

# Lock the audit configuration (must be last rule)
-e 2
