# files/custom_rules/0800-attack-detection.xml
<!--
  Wazuh Custom Attack Detection Rules
  Deployed by Ansible

  Rule ID range: 800100-800199 (avoids SOCFortress community rule ranges)
  These rules detect common attack patterns with MITRE ATT&CK mapping.
-->

<group name="custom_attack_detection,">

  <!-- в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ
       Reverse Shell Detection (T1059 - Command and Scripting Interpreter)
       в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ -->
  <rule id="800100" level="14">
    <if_sid>530,531</if_sid>
    <regex>bash -i >&amp; /dev/tcp|bash -i >&amp; /dev/udp|nc -e /bin|ncat -e /bin|python.*socket.*connect|perl.*socket.*INET|ruby.*TCPSocket|php.*fsockopen.*exec</regex>
    <description>Possible reverse shell detected: suspicious network-capable command execution</description>
    <mitre>
      <id>T1059</id>
      <id>T1071.001</id>
    </mitre>
    <group>attack,reverse_shell,</group>
  </rule>

  <rule id="800101" level="14">
    <if_sid>530,531</if_sid>
    <regex>mkfifo.*nc |mknod.*nc |/dev/tcp/|socat.*EXEC:|exec \d+&lt;&gt;/dev/tcp/</regex>
    <description>Possible reverse shell detected: named pipe or /dev/tcp redirection</description>
    <mitre>
      <id>T1059</id>
      <id>T1071.001</id>
    </mitre>
    <group>attack,reverse_shell,</group>
  </rule>

  <!-- в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ
       Credential Dumping (T1003 - OS Credential Dumping)
       в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ -->
  <rule id="800110" level="14">
    <if_sid>530,531</if_sid>
    <regex>mimikatz|sekurlsa|lsadump|pypykatz|gsecdump|wce.exe|lazagne|hashdump</regex>
    <description>Credential dumping tool detected</description>
    <mitre>
      <id>T1003</id>
    </mitre>
    <group>attack,credential_access,</group>
  </rule>

  <rule id="800111" level="12">
    <if_sid>530,531</if_sid>
    <regex>unshadow |john |hashcat |hydra |medusa |ncrack |patator</regex>
    <description>Password cracking/brute-force tool execution detected</description>
    <mitre>
      <id>T1110</id>
    </mitre>
    <group>attack,credential_access,</group>
  </rule>

  <rule id="800112" level="12">
    <if_sid>530,531</if_sid>
    <regex>/etc/shadow|/etc/gshadow|SAM.*SYSTEM.*SECURITY|secretsdump</regex>
    <description>Possible credential file access or extraction</description>
    <mitre>
      <id>T1003.008</id>
    </mitre>
    <group>attack,credential_access,</group>
  </rule>

  <!-- в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ
       Container Escape (T1611 - Escape to Host)
       в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ -->
  <rule id="800120" level="14">
    <if_sid>530,531</if_sid>
    <regex>nsenter.*-t.*1|docker.*--privileged|mount.*/proc/.*ns/|cdk.*exploit|/var/run/docker.sock</regex>
    <description>Possible container escape attempt detected</description>
    <mitre>
      <id>T1611</id>
    </mitre>
    <group>attack,container_escape,</group>
  </rule>

  <rule id="800121" level="12">
    <if_sid>530,531</if_sid>
    <regex>capsh.*--print|cat.*/proc/1/cgroup|mount.*cgroup|/proc/sysrq-trigger</regex>
    <description>Container environment reconnaissance or breakout attempt</description>
    <mitre>
      <id>T1611</id>
      <id>T1082</id>
    </mitre>
    <group>attack,container_escape,</group>
  </rule>

  <!-- в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ
       SSH Tunneling / Port Forwarding (T1572 - Protocol Tunneling)
       в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ -->
  <rule id="800130" level="10">
    <if_sid>530,531</if_sid>
    <regex>ssh.*-L \d+:|ssh.*-R \d+:|ssh.*-D \d+|ssh.*-w \d+:\d+|sshuttle|chisel </regex>
    <description>SSH tunnel or port forwarding detected</description>
    <mitre>
      <id>T1572</id>
    </mitre>
    <group>attack,lateral_movement,</group>
  </rule>

  <!-- в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ
       Ransomware Indicators (T1486 - Data Encrypted for Impact)
       в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ -->
  <rule id="800140" level="14">
    <if_sid>550,553,554</if_sid>
    <regex>.encrypted$|.locked$|.crypted$|.crypt$|.enc$|.ransom$|README.*DECRYPT|RECOVER.*FILES|HOW.*UNLOCK</regex>
    <description>Possible ransomware activity: mass file encryption or ransom note detected</description>
    <mitre>
      <id>T1486</id>
    </mitre>
    <group>attack,ransomware,</group>
  </rule>

  <!-- в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ
       Privilege Escalation via SUID/SGID (T1548.001)
       в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ -->
  <rule id="800150" level="10">
    <if_group>syscheck</if_group>
    <regex>Permissions changed.*suid|Permissions changed.*sgid</regex>
    <field name="changed_attributes">permission</field>
    <description>SUID/SGID bit set on file - possible privilege escalation setup</description>
    <mitre>
      <id>T1548.001</id>
    </mitre>
    <group>attack,privilege_escalation,</group>
  </rule>

  <!-- в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ
       Suspicious Cron Modification (T1053.003 - Cron)
       в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ -->
  <rule id="800151" level="10">
    <if_group>syscheck</if_group>
    <match>/etc/cron|/var/spool/cron</match>
    <description>Cron job modification detected - review for persistence</description>
    <mitre>
      <id>T1053.003</id>
    </mitre>
    <group>attack,persistence,</group>
  </rule>

  <!-- в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ
       Systemd Service Persistence (T1543.002)
       в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ -->
  <rule id="800152" level="10">
    <if_group>syscheck</if_group>
    <match>/etc/systemd/system|/usr/lib/systemd/system</match>
    <regex>.service$|.timer$</regex>
    <description>Systemd service/timer file modified - review for persistence</description>
    <mitre>
      <id>T1543.002</id>
    </mitre>
    <group>attack,persistence,</group>
  </rule>

  <!-- в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ
       Authorized Keys Modification (T1098.004)
       в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ -->
  <rule id="800153" level="10">
    <if_group>syscheck</if_group>
    <match>authorized_keys</match>
    <description>SSH authorized_keys file modified - review for unauthorized access persistence</description>
    <mitre>
      <id>T1098.004</id>
    </mitre>
    <group>attack,persistence,</group>
  </rule>

  <!-- в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ
       Reconnaissance Commands (T1082, T1016, T1049)
       в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ -->
  <rule id="800160" level="6">
    <if_sid>530,531</if_sid>
    <regex>whoami|id -a|uname -a|cat /etc/os-release|hostnamectl|systemctl list-units</regex>
    <description>System reconnaissance command executed</description>
    <mitre>
      <id>T1082</id>
    </mitre>
    <group>recon,discovery,</group>
  </rule>

  <rule id="800161" level="8">
    <if_sid>530,531</if_sid>
    <regex>nmap |masscan |zmap |rustscan |enum4linux|nbtscan|smbclient.*-L</regex>
    <description>Network scanning tool execution detected</description>
    <mitre>
      <id>T1046</id>
    </mitre>
    <group>attack,discovery,</group>
  </rule>

  <!-- в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ
       Defense Evasion (T1070 - Indicator Removal)
       в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ -->
  <rule id="800170" level="12">
    <if_sid>530,531</if_sid>
    <regex>history -c|history -w /dev/null|export HISTSIZE=0|unset HISTFILE|rm.*\.bash_history|shred.*\.bash_history</regex>
    <description>Shell history clearing detected - possible evidence tampering</description>
    <mitre>
      <id>T1070.003</id>
    </mitre>
    <group>attack,defense_evasion,</group>
  </rule>

  <rule id="800171" level="12">
    <if_sid>530,531</if_sid>
    <regex>rm -rf /var/log|truncate.*-s 0.*/var/log|> /var/log/|shred.*/var/log/</regex>
    <description>Log file deletion or truncation detected - possible evidence tampering</description>
    <mitre>
      <id>T1070.002</id>
    </mitre>
    <group>attack,defense_evasion,</group>
  </rule>

  <!-- в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ
       Kernel Module Loading (T1547.006)
       в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ -->
  <rule id="800180" level="10">
    <if_sid>80700,80701,80702</if_sid>
    <field name="audit.key">kernel_module</field>
    <description>Kernel module operation detected via auditd</description>
    <mitre>
      <id>T1547.006</id>
    </mitre>
    <group>attack,persistence,</group>
  </rule>

  <!-- в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ
       LD_PRELOAD Hijacking (T1574.006)
       в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ -->
  <rule id="800181" level="12">
    <if_group>syscheck</if_group>
    <match>/etc/ld.so.preload</match>
    <description>ld.so.preload modified - possible library injection persistence</description>
    <mitre>
      <id>T1574.006</id>
    </mitre>
    <group>attack,persistence,privilege_escalation,</group>
  </rule>

  <!-- в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ
       Data Exfiltration Indicators (T1048)
       в•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђв•ђ -->
  <rule id="800190" level="10">
    <if_sid>530,531</if_sid>
    <regex>curl.*-T |curl.*--upload|wget.*--post-file|scp.*@.*:|rsync.*@.*:</regex>
    <description>Possible data exfiltration: file upload to external destination</description>
    <mitre>
      <id>T1048</id>
    </mitre>
    <group>attack,exfiltration,</group>
  </rule>

</group>
